# cron/fetch.cron
# Executes *4* python (can be shell as well) scripts within a virtual environment sequentially, `fetch_active` -> `fetch_completed` -> `populate_active` -> `populate_completed` every day at 5:01 UTC (12:01 CST).
1 5 * * * /home/coopes/abupower/abupowerenv/bin/python /home/coopes/abupower/fetch_active_products.py >> /home/coopes/abupower/cronlog.txt ; /home/coopes/abupower/abupowerenv/bin/python /home/coopes/abupower/fetch_completed_products.py >> /home/coopes/abupower/cronlog.txt ; /home/coopes/abupower/abupowerenv/bin/python /home/coopes/abupower/populate_active_meta_data.py >> /home/coopes/abupower/cronlog.txt ; /home/coopes/abupower/abupowerenv/bin/python /home/coopes/abupower/populate_completed_meta_data.py >> /home/coopes/abupower/cronlog.txt 2>&1
